# mocha.service
[Unit]
Description=Mocha Cat Terminal
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/mocha
Environment=OLLAMA_HOST=http://localhost:11434
Environment=OLLAMA_MODEL=qwen3.5:9b
ExecStart=/opt/mocha/venv/bin/uvicorn app:app --host 0.0.0.0 --port 8400
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
